# config/supervisord.ini
[program:kibana]
command=/opt/kibana-4.2.0-linux-x64/bin/kibana
directory=/home/runtime/run/
autostart=true
autorestart=true
user=runtime
stdout_logfile=/home/runtime/run/kibana.log
stderr_logfile=/home/runtime/run/kibana-errors.log
